# infra/.service
[Unit]
Description=Coder gem API
After=network.target

[Service]
Type=simple
# Another Type: forking
User=nanodano
WorkingDirectory=/var/apps/coder-gem-api/coder-gem-api.jar

ExecStart=java -jar coder-gem-api.jar
Restart=on-failure
ExecStop=killall coder-gem-api
# Other restart options: always, on-abort, etc

# The install section is needed to use
# `systemctl enable` to start on boot
# For a user service that you want to enable
# and start automatically, use `default.target`
# For system level services, use `multi-user.target`
[Install]
WantedBy=multi-user.target
